# kde/personal/home/.config/autostart/launch.sh.desktop
[Desktop Entry]
Exec=/home/eric/.config/launch.sh
Icon=dialog-scripts
Name=launch.sh
Path=
Type=Application
X-KDE-AutostartScript=true
