# kde/personal/home/.config/autostart/vmware-user.desktop
[Desktop Entry]
Exec=vmware-user
Icon=
Name=vmware-user
Path=
Terminal=False
Type=Application