# kde/personal/home/.config/autostart/vmware-user-suid-wrapper.desktop
[Desktop Entry]
Exec=/usr/bin/vmware-user-suid-wrapper
Icon=
Name=vmware-user-suid-wrapper
Path=
Terminal=False
Type=Application